# scripts/monitoring/B2X-heartbeat.service
[Unit]
Description=B2X Runtime Health Check Heartbeat
After=network.target

[Service]
Type=simple
User=B2X
WorkingDirectory=/path/to/B2X
ExecStart=/path/to/B2X/scripts/runtime-health-check.sh --heartbeat --slack-webhook https://hooks.slack.com/services/YOUR/SLACK/WEBHOOK
Restart=always
RestartSec=5
Environment=SLACK_WEBHOOK_URL=https://hooks.slack.com/services/YOUR/SLACK/WEBHOOK

[Install]
WantedBy=multi-user.target